# infra/pi-network/systemd/smartcane-network-bootstrap.service
[Unit]
Description=SmartCane Network Bootstrap
Before=network-online.target smartcane-runtime.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=__SMARTCANE_REPO_ROOT__/infra/pi-network/bootstrap_network_mode.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
