# infra/pi-network/systemd/smartcane-network-bootstrap.service
[Unit]
Description=SmartCane Network Bootstrap (Deprecated - handled by main.py)
Before=network.target

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
