# infra/pi-network/systemd/smartcane-watchdog.service
[Unit]
Description=SmartCane Network Watchdog
After=network.target

[Service]
Type=simple
ExecStart=/bin/bash -c 'while true; do \
  if ! systemctl is-active --quiet hostapd; then systemctl restart hostapd; fi; \
  if ! systemctl is-active --quiet dnsmasq; then systemctl restart dnsmasq; fi; \
  sleep 30; \
done'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
