# infra/pi-network/systemd/smartcane-runtime.service
[Unit]
Description=SmartCane Pi Runtime
After=network-online.target smartcane-network-bootstrap.service
Wants=network-online.target smartcane-network-bootstrap.service

[Service]
Type=simple
WorkingDirectory=__SMARTCANE_REPO_ROOT__/pi
ExecStart=__SMARTCANE_REPO_ROOT__/pi/.venv/bin/python __SMARTCANE_REPO_ROOT__/pi/src/main.py
Restart=always
RestartSec=2
User=pi
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
