# pi/infra/pi-network/systemd/smartcane-runtime.service
[Unit]
Description=SmartCane Pi Runtime
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=__SMARTCANE_REPO_ROOT__/runtime
ExecStart=__SMARTCANE_REPO_ROOT__/infra/pi-network/run_pi_runtime.sh
ExecReload=/bin/sh -c "kill -HUP \$(cat /run/smartcane-runtime.pid)"
Restart=always
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1
Environment=SMARTCANE_LOG_LEVEL=INFO
WatchdogSec=300
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
